# fs/etc/systemd/system/sshd.service.d/10-require-mount-host-keys.conf
[Unit]
RequiresMountsFor=/etc/dracut-sshd/host_keys
