# fs/etc/systemd/system/sshd.service.d/10-require-mount-host-keys.conf
[Unit]
RequiresMountsFor=/var/lib/dracut-sshd
