# fs/usr/local/lib/systemd/system/dracut-sshd-copy-keys.path
[Path]
PathChanged=/root/.ssh/
PathChanged=/etc/ssh/

[Install]
WantedBy=multi-user.target
