# fs/etc/systemd/system/sshd-keygen@.service.d/10-copy.conf
[Service]
ExecStart=-/usr/sbin/sh -c "cp /etc/ssh/ssh_host_%i_key{,.pub} /etc/dracut-sshd/host_keys"
